# crontab
# sudo
0 0 1 */1 * certbot renew --renew-hook "systemctl reload nginx"

# user
*/15 * * * * python3 subreply/manage.py email
*/20 * * * * python3 subreply/manage.py fetch

30 0 * * * pg_dump -U postgres subreply | xz > db/subreply.sql.xz

45 */4 * * * python3 logparser/parse.py logs/sub.log.gz --html subreply/static/logs.html --ignore "subreply.com" --hide 1
