# crontab
# sudo
0 0 1 */1 * certbot renew --renew-hook "systemctl reload nginx"

# user
*/20 * * * * python3 ~/newscafe/manage.py fetch
30 12 * * * pg_dump -U postgres subreply | xz > ~/db/subreply.sql.xz
15 * * * * zcat -f ~/logs/news.log.gz | goaccess --ignore-crawlers --output ~/newscafe/static/logs.html -
45 * * * * zcat -f ~/logs/sub.log.gz | goaccess --ignore-crawlers --output ~/subreply/static/logs.html -
